# recipes-services/esync-k3s-deployment/files/esync-deployment-helper.service.in
[Unit]
Description=Deploy eSync-Client Demo to SOAFEE Framework (k3s)
Wants=esync-macaddr-monitor.service k3s.service
After=esync-macaddr-monitor.service k3s.service

[Service]
Type=oneshot
ExecStart=@ESCLIENT_PRIV_DIR@/kubeconfigs/deploy.sh
ExecStartPost=/bin/systemctl disable esync-deployment-helper.service

[Install]
WantedBy=multi-user.target
